# templates/default/mithril-service.conf.erb
# Dropped off by Chef recipe[mithril::cluster]
description "Mithril Service"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

setuid mithril
umask 022

respawn
respawn limit 3 10
chdir <%= @mithril_home %>
script
  <%= @mithril_home %>/bin/mithril-server \
    <% if @pg_enabled -%>
    -s postgresql
    -u $(awk -F: '{ print "postgres://" $4 ":" $5 "@" $1 "/" $3 "?sslmode=disable"}' <%= @mithril_home %>/.pgpass 2>/dev/null || echo '') \
    <% end -%>
    <% if @debug_enabled -%>
    -d \
    <% end -%>
    -p '<%= @pid_file %>' \
    '<%= @server_address %>' \
    '<%= @amqp_uri %>'
end script
